# handbrake/src/service/guacd/run
#!/bin/bash
exec 2>&1
exec /usr/local/sbin/guacd -f
